# pfv/pfv_worker/implementations/quic-implementations/quant/misc/quant.service.in
# systemd unit file to keep quant server running

[Unit]
Description=Quant server
StartLimitIntervalSec=0

[Service]
ExecStart=@CMAKE_CURRENT_BINARY_DIR@/quant-server.sh
Restart=always
RestartSec=1s

[Install]
WantedBy=default.target
